# config/launchd/com.local-inference.terminal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local-inference.terminal</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/ttyd</string>
        <string>-p</string>
        <string>7681</string>
        <string>-i</string>
        <string>127.0.0.1</string>
        <string>-W</string>
        <string>zsh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/local-inference/terminal.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/local-inference/terminal.err.log</string>
</dict>
</plist>
